# linux-voice.service
[Unit]
Description=Voice-to-text dictation tool
Documentation=https://github.com/williamjameshandley/linux-voice
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/linux-voice
Restart=on-failure
RestartSec=5

# Inherit DISPLAY from graphical-session environment

[Install]
WantedBy=graphical-session.target
